# Scripts/clightd.service
[Unit]
Description=Bus service to manage various screen related properties (gamma, dpms, backlight)
Requires=@POLKIT_NAME@.service

[Service]
Type=dbus
BusName=org.clightd.clightd
User=root
ExecStart=@DAEMON_DIR@/clightd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
